# rootfs/etc/s6-overlay/s6-rc.d/init-nodered/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Node-RED
# Configure Node-RED avant exécution
# ==============================================================================

# Assure que la configuration existe
if ! bashio::fs.directory_exists '/config/smartphoton_jkbms/'; then
    mkdir -p /config/smartphoton_jkbms/nodes \
        || bashio::exit.nok "Failed to create node-red configuration directory"

    # Copie les fichiers modèles
    cp /etc/node-red/flows.json /config/smartphoton_jkbms/
    cp /etc/node-red/settings.js /config/smartphoton_jkbms/

    # Création d'un ID de flow aléatoire
    id=$(node -e "console.log((1+Math.random()*4294967295).toString(16));")
    sed -i "s/%%ID%%/${id}/" "/config/smartphoton_jkbms/flows.json"
fi

# Vérification de version
if ! bashio::fs.file_exists '/config/smartphoton_jkbms/version'; then
    touch /config/smartphoton_jkbms/version
    VERSION=$(bashio::addon.version)
    bashio::log.info "Module Addon Version number: ${VERSION}"
    echo $VERSION >> /config/smartphoton_jkbms/version
    chmod 775 /config/smartphoton_jkbms/version
    bashio::log.info "${VERSION}"
else
    VERSIONOLD=$(cat /config/smartphoton_jkbms/version)
    VERSION=$(bashio::addon.version)
    if [ $VERSIONOLD != $VERSION ]; then
        bashio::log.info "Mise a jour de la nouvelle version"
        bashio::log.info "Module Addon Version number: ${VERSION}"
        if bashio::fs.file_exists '/config/smartphoton_jkbms/flows.backup'; then
            rm /config/smartphoton_jkbms/flows.backup
        fi
        mv /config/smartphoton_jkbms/flows.json /config/smartphoton_jkbms/flows.backup
        cp /etc/node-red/flows.json /config/smartphoton_jkbms/
    fi
fi

# Remplace le fichier flows par celui d'origine
if bashio::fs.file_exists '/config/smartphoton_jkbms/flows.backup'; then
    rm /config/smartphoton_jkbms/flows.backup
fi
mv /config/smartphoton_jkbms/flows.json /config/smartphoton_jkbms/flows.backup
cp /etc/node-red/flows.json /config/smartphoton_jkbms/

# Désinstallation des packages Node-RED conflictuels
cd /config || bashio::exit.nok "Could not change directory to Node-RED"
if bashio::fs.file_exists "/config/package.json"; then
    npm uninstall \
        node-red-contrib-home-assistant \
        node-red-contrib-home-assistant-llat \
        node-red-contrib-home-assistant-ws \
        || bashio::log.warning "Failed un-installing conflicting packages"
fi
